# awtrix-pomo.service
[Unit]
Description=Awtrix Pomodoro web trigger
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/python3 %h/git/mqtt/awtrix_pomo_server.py
Restart=on-failure
RestartSec=3

# Hardening — safe because the script only reads its own code,
# opens an inbound TCP socket, and talks to the local MQTT broker.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=default.target
